# packages/rpicammqtt/rpicammqtt-0.10.1-py3-none-any.whl/rpicammqtt/systemd/rpi-cam-mqtt.service
[Unit]
Description=rpi-cam-mqtt
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rpicammqtt

[Install]
WantedBy=multi-user.target
